# roles/newebe/files/newebe.conf
[program:newebe]
autorestart=false
command=newebe_server.py --configfile=/usr/local/etc/newebe/config.yaml
redirect_stderr=true
user=newebe
