# roles/mailserver/files/lib_systemd_system_rmilter.socket
[Unit]
Description=Another sendmail milter for different mail checks

[Socket]
ListenStream=9900

[Install]
WantedBy=sockets.target
